# src/etc/systemd/system/sssd.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/ip vrf exec mgmt /usr/sbin/sssd -i ${DEBUG_LOGGER}
LimitMEMLOCK=8388608
