# src/etc/systemd/system/nslcd.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/ip vrf exec mgmt /etc/init.d/nslcd start
